# server/deploy/axdea.service
# systemd 유닛 (pm2 대신 쓰고 싶을 때). 경로/사용자는 실제 환경에 맞게 수정하세요.
# 설치:
#   sudo cp axdea.service /etc/systemd/system/axdea.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now axdea
#   journalctl -u axdea -f      # 로그
[Unit]
Description=AXdea web server (Express + MySQL)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# ↓ 실제 배포 경로/사용자로 변경
WorkingDirectory=/opt/axdea/server
ExecStart=/usr/bin/node server.js
EnvironmentFile=/opt/axdea/server/.env
Restart=always
RestartSec=3
User=axdea
Group=axdea

[Install]
WantedBy=multi-user.target
